# dot_config/systemd/user/kds-docs-borg-s3-backup.timer
[Unit]
Description=Back up ~/Documents to S3 daily

[Timer]
OnCalendar=*-*-* 17:15
Persistent=true

[Install]
WantedBy=timers.target
