# dot_config/systemd/user/kds-docs-borg-s3-backup.timer
[Unit]
Description=Daily backup ~/Documents to S3

[Timer]
OnCalendar=*-*-* 17:15
Persistent=true

[Install]
WantedBy=timers.target
